# 05a - Post YML Docker Fix
################# Post YML Fix ################# START
#### Problem:  Couch and Sonnar are unable to see the Google Drive via RCLONE
#### Reason:  RClone loads up mores slowly and the docker containers load too fast
#### Answer:  Force a reboot of those two containers after 25 seconds
#### Result:  This script reboots both containers; now they should see the RClone directory

### location of of postyml script
sudo nano /opt/postyml.sh

##### Copy #### START
#!/bin/bash

### BEGIN INIT INFO
# Provides:          postyml.sh
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: postyml.sh
# Description:       postyml-fix
### END INIT INFO

x=20
while [ $x -gt 0 ]
do
    sleep 1s
    clear
    echo "$x seconds until reboot"
    x=$(( $x - 1 ))
done
docker restart couchpotato
docker restart sonarr
docker restart plex

exit 0;
##### Copy #### END

# Press CTRL+X; then Enter. Editing rc.local to add postyml fix
sudo nano /etc/rc.local

# Place the line below above "exit 0" | If your missing the rc.local script and it shows new; then goto the very bottom
sh /opt/postyml.sh &
# Press CTRL+X and then ENTER to save

# Assign permissions so your script can execute
sudo chown root /etc/rc.local
sudo chmod 755 /etc/rc.local
sudo chmod 755 /opt/postyml.sh
sudo chown root /opt/postyml.sh
sudo /etc/init.d/rc.local start
################# Post YML FIX ################# END

################# Only for Missing rc.local ###################### START
########################### Make the script look like this, then resume above at # Assign permissions so ...


#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

sh /opt/postyml.sh &

exit 0
################# Only for Missing rc.local ###################### END
