# 05a - Post YML Docker Fix
################# Post YML Fix ################# START
#### Problem:  Couch and Sonnar are unable to see the Google Drive via RCLONE
#### Reason:  RClone loads up mores slowly and the docker containers load too fast
#### Answer:  Force a reboot of those two containers after 25 seconds
#### Result:  This script reboots both containers; now they should see the RClone directory

### location of of postyml script
sudo nano /opt/postyml.sh

##### Copy #### START
#!/bin/bash

### BEGIN INIT INFO
# Provides:          postyml.sh
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: postyml.sh
# Description:       postyml-fix
### END INIT INFO

x=20
while [ $x -gt 0 ]
do
    sleep 1s
    clear
    echo "$x seconds until reboot"
    x=$(( $x - 1 ))
done
docker restart couchpotato
docker restart sonarr
docker restart plex

exit 0;
##### Copy #### END

# Press CTRL+X; then Enter.

sudo chown root:root /opt/postyml.sh
sudo chmod 755 /opt/postyml.sh

sudo nano /etc/systemd/system/postyml.service

##### START COPY #####
[Unit]
Description=rclonemount
After=multi-user.target

[Service]
ExecStart=/bin/bash /opt/postyml.sh
ExecStop=/bin/bash /opt/postyml.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target

##### END COPY #####

sudo systemctl daemon-reload
sudo systemctl enable postyml.service
sudo systemctl start postyml.service
sudo systemctl status postyml.service

################# POST YML FIX ###################### END
