# 05b - Post YML Docker Fix
################# Post YML Fix ################# START
#### Problem:  Couch and Sonnar are unable to see the Google Drive via RCLONE
#### Reason:  RClone loads up mores slowly and the docker containers load too fast
#### Answer:  Force a reboot of those two containers after 25 seconds
#### Result:  They reboot and both containers should see the RCLONE Script
y
### PlexDrive - RC.Local startup script
sudo nano /opt/postyml.sh

##### Copy #### START
#!/bin/bash

### BEGIN INIT INFO
# Provides:          postyml.sh
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: postyml.sh
# Description:       postyml-fix
### END INIT INFO

x=25
while [ $x -gt 0 ]
do
    sleep 1s
    clear
    echo "$x seconds until blast off"
    x=$(( $x - 1 ))
done
docker restart couchpotato
docker restart sonar

exit 0;
##### Copy #### END

# Press CTRL+X; then Enter. Editing rc.local to add postyml fix
sudo nano /etc/rc.local

# Place the line below above "exit 0"
sh /opt/plex5.sh & 
sh /opt/postyml.sh &
# Press CTRL+X and then ENTER to save

# Assign permissions so you script can execute
sudo chown root /etc/rc.local
sudo chmod 755 /etc/rc.local
sudo chmod 755 /opt/postyml.sh
sudo chown root /opt/postyml.sh
sudo /etc/init.d/rc.local start
################# Post YML FIX ################# END
