# ansible/roles/plexdrive/templates/plexdrive.service
[Unit]
Description=PlexDrive Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/plexdrive mount --uid=1000 --gid=1000 --uid=1000 --gid=1000 -v 3 --refresh-interval=1m --chunk-check-threads=8 --chunk-load-threads=8 --chunk-load-ahead=6 --fuse-options=allow_other,read_only,allow_non_empty_mount --config=/root/.plexdrive --cache-file=/root/.plexdrive/cache.bolt /mnt/plexdrive
ExecStop=/bin/fusermount -uz /mnt/plexdrive
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target